# scripts/cropnet.service
[Unit]
Description=CropNet API — Multi-Model Crop Disease Diagnosis
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/cropnet
EnvironmentFile=/opt/cropnet/.env
Environment=IMG_SIZE=256
ExecStart=/opt/cropnet/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
